# supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisor/supervisord.pid
nodaemon=true
user=root


[program:postfix]
command=/usr/sbin/postfix start-fg

